# systemd/dnspulse.service
# /etc/systemd/system/dnspulse.service
[Unit]
Description=DNS Pulse Exporter - Prometheus DNS monitoring service
Documentation=https://github.com/farrokhi/dnspulse_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/dnspulse_exporter -f /etc/dnspulse.yml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
TimeoutStopSec=20

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"

# Security Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6

# Resource limits
TasksMax=10
LimitNOFILE=1024
LimitNPROC=10

# Run as daemon user
User=daemon
Group=daemon

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dnspulse_exporter

[Install]
WantedBy=multi-user.target
